# .local/runit/sv/ipfs/run
#!/bin/bash

if [ ! -d $HOME/.ipfs ]; then
    ipfs init
fi

exec ipfs daemon
