# .local/runit/sv/mpd/run
#!/bin/sh

# Wait for pipewire-pulse
sv check "$svdir/pipewire-pulse" >/dev/null || exit 1

[ -r ./conf ] && . ./conf

exec mpd --no-daemon --stderr ${OPTS:-} 2>&1
